# systemd/calendash.service
[Unit]
Description=Generate daily calendar dashboard image
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/pihole/zero2dash
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=-/home/pihole/zero2dash/.env
TimeoutStartSec=2min
ExecStart=/usr/bin/python3 -u /home/pihole/zero2dash/modules/calendash/calendash-api.py

