# systemd/day.timer
[Unit]
Description=Start DAY GUI at 07:00

[Timer]
OnCalendar=*-*-* 07:00:00
AccuracySec=1min
Unit=display.service

[Install]
WantedBy=timers.target
