# systemd/display.service
[Unit]
Description=Pi-hole TFT Day Display Rotator
After=network-online.target
Wants=network-online.target
Conflicts=night.service
Before=night.service

[Service]
Type=simple
WorkingDirectory=/home/pihole/zero2dash
Environment=PYTHONUNBUFFERED=1
Environment=FB_DEVICE=/dev/fb1
EnvironmentFile=-/home/pihole/zero2dash/.env

ExecStartPre=/home/pihole/zero2dash/pihole-display-pre.sh
ExecStart=/usr/bin/python3 -u /home/pihole/zero2dash/display_rotator.py

Restart=always
RestartSec=2
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

