# systemd/night.service
[Unit]
Description=TFT Night Blackout Display
After=network-online.target
Wants=network-online.target
Conflicts=boot-selector.service
Conflicts=display.service
Before=display.service

[Service]
Type=simple
WorkingDirectory=/home/pihole/zero2dash
Environment=PYTHONUNBUFFERED=1
Environment=FB_DEVICE=/dev/fb1
EnvironmentFile=-/home/pihole/zero2dash/.env

ExecStartPre=/home/pihole/zero2dash/pihole-display-pre.sh
ExecStart=/usr/bin/python3 -u /home/pihole/zero2dash/modules/blackout/blackout.py

Restart=always
RestartSec=2
StandardOutput=journal
StandardError=journal

# Manual compatibility foreground path retained for legacy Night launches.

