# systemd/pihole-display.service
[Unit]
Description=Pi-hole TFT Stats Display
After=network-online.target dev-fb1.device
Wants=network-online.target dev-fb1.device
ConditionPathExists=/dev/fb1

[Service]
Type=simple
WorkingDirectory=/home/pihole
Environment=PYTHONUNBUFFERED=1

ExecStartPre=/usr/local/bin/pihole-display-pre.sh
ExecStart=/usr/bin/python3 -u /home/pihole/piholestats_v1.1.py

Restart=always
RestartSec=2
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
