# debug/currency-salvage-2026-03-09/feature-branch/systemd/currency-update.service
[Unit]
Description=Refresh GBP/PLN currency image
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/zero2dash
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=-/opt/zero2dash/.env
ExecStart=/usr/bin/python3 -u /opt/zero2dash/scripts/currency-rate.py

[Install]
WantedBy=multi-user.target
